# src/conf/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

;[inet_http_server]         ; inet (TCP) server disabled by default
;port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
;umask=022                   ; (process file creation umask;default 022)
;user=chrism                 ; (default is current user, required if root)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=/tmp              ; (default is not to cd during start)
;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;environment=KEY="value"     ; (key value pairs to add to environment)
;strip_ansi=false            ; (strip ansi escape codes in logs; def. false)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

; ==================================
; app-server superviord
; ==================================

[program:marsweb8890]
command=python /search/index.py --port=8890
directory=/search/
user=root
numprocs=1
redirect_stderr=true        ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/marsweb_std.log
stderr_logfile=/var/log/marsweb_error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true
priority=998
stdout_logfile_maxbytes = 20MB
stdoiut_logfile_backups = 20

[program:marsweb8891]
command=python /search/index.py --port=8891
directory=/search/
user=root
numprocs=1
redirect_stderr=true        ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/marsweb_std.log
stderr_logfile=/var/log/marsweb_error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true
priority=998
stdout_logfile_maxbytes = 20MB
stdoiut_logfile_backups = 20

#[program:marsweb8892]
#command=python /search/index.py --port=8892
#directory=/search/
#user=root
#numprocs=1
#redirect_stderr=true        ; redirect proc stderr to stdout (default false)
#stdout_logfile=/var/log/marsweb_std.log
#stderr_logfile=/var/log/marsweb_error.log
#autostart=true
#autorestart=true
#startsecs=10
#stopwaitsecs = 600
#killasgroup=true
#priority=998
#stdout_logfile_maxbytes = 20MB
#stdoiut_logfile_backups = 20
#
#[program:marsweb8893]
#command=python /search/index.py --port=8893
#directory=/search/
#user=root
#numprocs=1
#redirect_stderr=true        ; redirect proc stderr to stdout (default false)
#stdout_logfile=/var/log/marsweb_std.log
#stderr_logfile=/var/log/marsweb_error.log
#autostart=true
#autorestart=true
#startsecs=10
#stopwaitsecs = 600
#killasgroup=true
#priority=998
#stdout_logfile_maxbytes = 20MB
#stdoiut_logfile_backups = 20
#
#[program:marsweb8894]
#command=python /search/index.py --port=8894
#directory=/search/
#user=root
#numprocs=1
#redirect_stderr=true        ; redirect proc stderr to stdout (default false)
#stdout_logfile=/var/log/marsweb_std.log
#stderr_logfile=/var/log/marsweb_error.log
#autostart=true
#autorestart=true
#startsecs=10
#stopwaitsecs = 600
#killasgroup=true
#priority=998
#stdout_logfile_maxbytes = 20MB
#stdoiut_logfile_backups = 20
#
#[program:marsweb8895]
#command=python /search/index.py --port=8895
#directory=/search/
#user=root
#numprocs=1
#redirect_stderr=true        ; redirect proc stderr to stdout (default false)
#stdout_logfile=/var/log/marsweb_std.log
#stderr_logfile=/var/log/marsweb_error.log
#autostart=true
#autorestart=true
#startsecs=10
#stopwaitsecs = 600
#killasgroup=true
#priority=998
#stdout_logfile_maxbytes = 20MB
#stdoiut_logfile_backups = 20
#
#[program:marsweb8896]
#command=python /search/index.py --port=8896
#directory=/search/
#user=root
#numprocs=1
#redirect_stderr=true        ; redirect proc stderr to stdout (default false)
#stdout_logfile=/var/log/marsweb_std.log
#stderr_logfile=/var/log/marsweb_error.log
#autostart=true
#autorestart=true
#startsecs=10
#stopwaitsecs = 600
#killasgroup=true
#priority=998
#stdout_logfile_maxbytes = 20MB
#stdoiut_logfile_backups = 20
#
#[program:marsweb8897]
#command=python /search/index.py --port=8897
#directory=/search/
#user=root
#numprocs=1
#redirect_stderr=true        ; redirect proc stderr to stdout (default false)
#stdout_logfile=/var/log/marsweb_std.log
#stderr_logfile=/var/log/marsweb_error.log
#autostart=true
#autorestart=true
#startsecs=10
#stopwaitsecs = 600
#killasgroup=true
#priority=998
#stdout_logfile_maxbytes = 20MB
#stdoiut_logfile_backups = 20
#
#[program:marsweb8898]
#command=python /search/index.py --port=8898
#directory=/search/
#user=root
#numprocs=1
#redirect_stderr=true        ; redirect proc stderr to stdout (default false)
#stdout_logfile=/var/log/marsweb_std.log
#stderr_logfile=/var/log/marsweb_error.log
#autostart=true
#autorestart=true
#startsecs=10
#stopwaitsecs = 600
#killasgroup=true
#priority=998
#stdout_logfile_maxbytes = 20MB
#stdoiut_logfile_backups = 20
#
#[program:marsweb8899]
#command=python /search/index.py --port=8899
#directory=/search/
#user=root
#numprocs=1
#redirect_stderr=true        ; redirect proc stderr to stdout (default false)
#stdout_logfile=/var/log/marsweb_std.log
#stderr_logfile=/var/log/marsweb_error.log
#autostart=true
#autorestart=true
#startsecs=10
#stopwaitsecs = 600
#killasgroup=true
#priority=998
#stdout_logfile_maxbytes = 20MB
#stdoiut_logfile_backups = 20